# image_scripts/rc.local
#!/bin/sh 
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.
/etc/rc.firstboot
bash /home/pi/nginx_config.sh > /home/nginx_config.log 2>&1 &
bash /home/pi/hotspot.sh > /home/hostpot.log 2>&1 &
/home/pi/status.sh > /home/status.log 2>&1 &
exit 0
